# ScholarScraper/crontab

0 12 * * * /usr/local/bin/bash /home/V00837208/ScholarScraper/ScholarScraper/run.sh
